# systemd/waydroid-cloud-phone.target
# Waydroid Cloud Phone - Main Target
# All services are grouped under this target for easy management

[Unit]
Description=Waydroid Cloud Phone Stack
Documentation=https://github.com/example/waydroid-cloud-phone
After=network-online.target
Wants=network-online.target

# Component services
Wants=nginx-rtmp.service
Wants=xvnc.service
Wants=waydroid-container.service
Wants=waydroid-session.service
Wants=ffmpeg-bridge.service
Wants=control-api.service
Wants=agent-api.service

[Install]
WantedBy=multi-user.target
