# systemd/waydroid-session.service
[Unit]
Description=Waydroid Session
Documentation=https://waydro.id
After=waydroid-container.service xvnc.service
Requires=waydroid-container.service
Wants=xvnc.service
ConditionPathExists=/var/lib/waydroid/waydroid.cfg

[Service]
Type=simple
User=waydroid
Group=waydroid

# Wait for container to be ready
ExecStartPre=/bin/sleep 5

ExecStart=/usr/bin/waydroid session start
ExecStop=/usr/bin/waydroid session stop

Restart=on-failure
RestartSec=10
TimeoutStartSec=120

# Environment for X11 (VNC)
Environment=DISPLAY=:1
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=XDG_SESSION_TYPE=x11
Environment=HOME=/home/waydroid
Environment="WAYLAND_DISPLAY="

[Install]
WantedBy=waydroid-cloud-phone.target
