# systemd/xvnc.service
[Unit]
Description=TigerVNC Server for Waydroid
Documentation=man:Xvnc(1)
After=network.target

[Service]
Type=simple
User=waydroid
Group=waydroid

# VNC on display :1 (port 5901)
# Geometry matches typical phone resolution
ExecStart=/usr/bin/Xvnc :1 \
    -geometry 1080x1920 \
    -depth 24 \
    -rfbauth /home/waydroid/.vnc/passwd \
    -rfbport 5901 \
    -localhost \
    -alwaysshared \
    -dpi 320 \
    -SecurityTypes VncAuth

Restart=on-failure
RestartSec=5

# Environment
Environment=HOME=/home/waydroid
Environment=USER=waydroid

[Install]
WantedBy=waydroid-cloud-phone.target
